# systemd/secret-broker-backup.service
[Unit]
Description=Run Secret Broker SQLite backup

[Service]
Type=oneshot
WorkingDirectory=/opt/secret-broker
ExecStart=/opt/secret-broker/scripts/backup-sqlite.sh /var/lib/secret-broker/secret-broker.db /var/backups/secret-broker/secret-broker.db.%Y%m%d%H%M%S.bak
